# ninjapack/root/etc/init/sphere-fix-applicator.conf
start on runlevel [23]
console log
task

env LOG="/var/log/sphere-fix-applicator.log"

script
	bash -c "(
		sleep 5
		apt-get remove ninja-fix-applicato
	) & disown" > $LOG 2>&1
end script
